# target/init_scripts/systemd-capemgr.service
[Unit]
Description=Cape Manager Service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/sh /opt/scripts/boot/capemgr.sh

[Install]
WantedBy=multi-user.target
