# target/init_scripts/systemd-generic-board-startup.service
[Unit]
Description=Generic Board Startup
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/sh /opt/scripts/boot/generic-startup.sh

[Install]
WantedBy=multi-user.target
